# meta-isar/recipes-graphics/weston-dk2-config/files/weston.service
[Unit]
Description=Weston Wayland Compositor
# Wait for the DRM device to appear via udev (handles deferred probe on LTDC).
# ConditionPathExists is evaluated once at start-up and does not retry on
# condition failure, so we bind to the device unit instead.
After=dev-dri-card0.device
Requires=dev-dri-card0.device

[Service]
Environment=XDG_RUNTIME_DIR=/run/user/0
ExecStartPre=/bin/mkdir -p /run/user/0
ExecStartPre=/bin/chmod 700 /run/user/0
ExecStart=/usr/bin/weston --backend=drm-backend.so --tty=2 --use-pixman
Restart=on-failure
StandardInput=tty
TTYPath=/dev/tty2

[Install]
WantedBy=multi-user.target
